# deploy/fiscalbay-alertcheck.service
[Unit]
Description=FiscalBay runtime alert check

[Service]
Type=oneshot
User=__APP_USER__
Group=__APP_GROUP__
WorkingDirectory=__APP_DIR__
ExecStart=/usr/bin/bash -lc 'exec /usr/bin/bash "__APP_DIR__/deploy/alert-check.sh"'
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=__APP_DIR__/data
MemoryMax=__ONESHOT_MEMORY_MAX__
CPUQuota=__ONESHOT_CPU_QUOTA__
TasksMax=__ONESHOT_TASKS_MAX__
