# deploy/fiscalbay-autodeploy.timer
[Unit]
Description=Auto-deploy timer for FiscalBay

[Timer]
OnBootSec=5min
OnUnitActiveSec=10min
Persistent=true
Unit=fiscalbay-autodeploy.service

[Install]
WantedBy=timers.target
