# deploy/fiscalbay-release-please.timer
[Unit]
Description=Run FiscalBay release-please Release PR automation hourly

[Timer]
OnCalendar=hourly
RandomizedDelaySec=15m
Persistent=true
Unit=fiscalbay-release-please.service

[Install]
WantedBy=timers.target
