# deploy/fiscalbay-alertcheck.timer
[Unit]
Description=Periodic alert check for FiscalBay

[Timer]
OnCalendar=*:0/5
Persistent=true
Unit=fiscalbay-alertcheck.service

[Install]
WantedBy=timers.target
