# deploy/fiscalbay-reconcile.service
[Unit]
Description=FiscalBay periodic reconciliation

[Service]
Type=oneshot
User=__APP_USER__
Group=__APP_GROUP__
WorkingDirectory=__APP_DIR__
ExecStart=/usr/bin/bash -c 'source "__VENV_DIR__/bin/activate" && set -a && source "__ENV_FILE__" && set +a && exec /usr/bin/bash "__APP_DIR__/deploy/reconcile.sh"'
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=__APP_DIR__/data
MemoryMax=__ONESHOT_MEMORY_MAX__
CPUQuota=__ONESHOT_CPU_QUOTA__
TasksMax=__ONESHOT_TASKS_MAX__
