# deploy/fiscalbay-reconcile.timer
[Unit]
Description=Periodic reconciliation for FiscalBay

[Timer]
OnCalendar=*:0/10
Persistent=true
Unit=fiscalbay-reconcile.service

[Install]
WantedBy=timers.target
