# deploy/fiscalbay-backup.timer
[Unit]
Description=Daily backup for FiscalBay

[Timer]
OnCalendar=*-*-* 03:15:00
Persistent=true
Unit=fiscalbay-backup.service

[Install]
WantedBy=timers.target
