# deploy/fiscalbay-external-healthcheck.service
[Unit]
Description=FiscalBay external HTTPS healthcheck
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=__APP_USER__
Group=__APP_GROUP__
WorkingDirectory=__APP_DIR__
ExecStart=/usr/bin/bash -c 'exec "__APP_DIR__/deploy/external-healthcheck.sh"'
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
MemoryMax=__ONESHOT_MEMORY_MAX__
CPUQuota=__ONESHOT_CPU_QUOTA__
TasksMax=__ONESHOT_TASKS_MAX__
