# deploy/fiscalbay-duckdns.service
[Unit]
Description=FiscalBay Duck DNS update
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/fiscalbay/duckdns.env
ExecStart=/opt/fiscalbay/deploy/duckdns-update.sh

