# deploy/fiscalbay-release-please.service
[Unit]
Description=FiscalBay release-please Release PR automation
Wants=network-online.target
After=network-online.target
ConditionPathExists=__APP_DIR__/release-please-config.json
ConditionPathExists=__APP_DIR__/.release-please-manifest.json

[Service]
Type=oneshot
WorkingDirectory=__APP_DIR__
Environment=APP_DIR=__APP_DIR__
Environment=HOME=/root
EnvironmentFile=-/etc/fiscalbay/release-please.env
ExecStart=__APP_DIR__/deploy/release-please-pr.sh
NoNewPrivileges=true
PrivateTmp=true
