# deploy/fiscalbay-log-maintenance.service
[Unit]
Description=FiscalBay log maintenance

[Service]
Type=oneshot
WorkingDirectory=__APP_DIR__
ExecStart=/usr/bin/bash -lc 'set -a && source "__ENV_FILE__" && set +a && exec "__APP_DIR__/deploy/log-maintenance.sh"'
