# deploy/fiscalbay-restore-drill.service
[Unit]
Description=FiscalBay restore drill

[Service]
Type=oneshot
User=__APP_USER__
Group=__APP_GROUP__
WorkingDirectory=__APP_DIR__
ExecStart=/usr/bin/bash -c 'exec "__APP_DIR__/deploy/restore-drill.sh"'
NoNewPrivileges=true
PrivateTmp=true
MemoryMax=__ONESHOT_MEMORY_MAX__
CPUQuota=__ONESHOT_CPU_QUOTA__
TasksMax=__ONESHOT_TASKS_MAX__
